# config/etc/supervisor/conf.d/nginx.conf
[program:nginx]
command=/usr/sbin/nginx -c /etc/nginx/nginx.conf
stopsignal=QUIT
redirect_stderr=true
autorestart=true